# CODE/systemd/spotify-kiosk.service
[Unit]
Description=Spotify + LED kiosk web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sophie
WorkingDirectory=/home/sophie/CODE/spotify_kiosk
EnvironmentFile=/etc/default/spotify-kiosk
ExecStart=/usr/bin/python3 /home/sophie/CODE/spotify_kiosk/app.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
